# fertz.service

[Unit]
Description=FERTZ service

[Service]
Type=forking
TimeoutStartSec=20
ExecStart=//usr/local/python/aquariummonitor/fertz.py 8